# vype-dev-vulkan.service
[Unit]
Description=Vype speech-to-text keyboard (dev, Vulkan)
After=graphical-session.target default.target sound.target
Wants=graphical-session.target

[Service]
Type=simple
TimeoutStartSec=15min
Environment=DISPLAY=:0
ExecStartPre=/usr/bin/env cargo build -r -p vyped --no-default-features --features vulkan
ExecStart=/bin/sh -c 'export DISPLAY="${DISPLAY:-:0}"; if [ -z "$XAUTHORITY" ] || [ ! -r "$XAUTHORITY" ]; then for c in /tmp/xauth_* "$HOME/.Xauthority" "$XDG_RUNTIME_DIR/Xauthority"; do [ -r "$c" ] && export XAUTHORITY="$c" && break; done; fi; exec %h/d/vype/target/release/vyped -s large -k F9'
WorkingDirectory=%h/d/vype
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60

[Install]
WantedBy=default.target
